# .config/dev/supervisord.conf
[program:app]
directory=/srv/app/
command=gunicorn chat_example.app:app -c /srv/app/.config/dev/gunicorn/config.py
autostart=true
autorestart=true

[program:nginx]
command=nginx

